# k8s-cert-sync/templates/k8s-cert-sync.service.j2
[Unit]
Description=K8s Certificate Sync Service

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/copy_kube_config.sh

[Install]
WantedBy=multi-user.target
